# wttr.timer
[Unit]
Description=a timer for wttr service
Requires=wttr.service

[Timer]
Unit=wttr.service
OnCalendar=hourly

[Install]
WantedBy=timers.target
